# claude/launchd/bio.claude.api-registry-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>bio.claude.api-registry-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>$HOME/local/bin/cc-api-find --refresh --force >> /tmp/cc-api-registry.log 2>&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>17</integer>
  </dict>
  <key>StandardOutPath</key><string>/tmp/cc-api-registry.log</string>
  <key>StandardErrorPath</key><string>/tmp/cc-api-registry.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- fnm node lives under aliases/default/bin; resolve through that symlink
         rather than the per-shell fnm_multishells path which is not stable. -->
    <key>PATH</key><string>/Users/leonardofibonacci/.local/share/fnm/aliases/default/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
