# claude/launchd/bio.claude.self-improve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>bio.claude.self-improve</string>
  <key>ProgramArguments</key>
  <array><string>/bin/bash</string><string>-c</string><string>$HOME/.claude/scripts/self-improvement-digest.sh 7</string></array>
  <key>StartCalendarInterval</key>
  <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>20</integer></dict>
  <key>StandardOutPath</key><string>/tmp/claude-digest.log</string>
  <key>StandardErrorPath</key><string>/tmp/claude-digest.err</string>
  <key>EnvironmentVariables</key>
  <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string></dict>
</dict>
</plist>
