# claude/launchd/com.jonathan.fabric.self-heal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.jonathan.fabric.self-heal</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/leonardofibonacci/.claude/fabric/self_heal.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/leonardofibonacci/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/leonardofibonacci</string>
  </dict>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/leonardofibonacci/.claude/fabric/logs/self-heal.log</string>
  <key>StandardErrorPath</key><string>/Users/leonardofibonacci/.claude/fabric/logs/self-heal.log</string>
</dict>
</plist>
