# launchd/com.user.cc-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent template for cc-loop. NOT YET INSTALLED.

  To install (after authorizing the cc-loop runner):
    1. Substitute the prompt path for your actual prompt file.
    2. Copy to ~/Library/LaunchAgents/com.user.cc-loop.plist.
    3. launchctl load -w ~/Library/LaunchAgents/com.user.cc-loop.plist.

  To pause without unloading:
    touch ~/.cc-paused              # cc-loop polls this and sleeps
  Resume:
    rm ~/.cc-paused
  Stop launchd from respawning:
    launchctl unload ~/Library/LaunchAgents/com.user.cc-loop.plist

  KeepAlive uses PathState on ~/.cc-paused as a hard-stop: when present,
  launchd won't restart the process even if it exits. Combined with the
  in-loop poll, this gives a defense-in-depth pause path.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.cc-loop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/leonardofibonacci/dotfiles/bin/cc-loop</string>
    <string>/Users/leonardofibonacci/.claude/cc-loop-prompt.md</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/leonardofibonacci</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>PathState</key>
    <dict>
      <!-- Run only while ~/.cc-paused does NOT exist -->
      <key>/Users/leonardofibonacci/.cc-paused</key>
      <false/>
    </dict>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/leonardofibonacci/.claude/cc-loop.stdout</string>

  <key>StandardErrorPath</key>
  <string>/Users/leonardofibonacci/.claude/cc-loop.stderr</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/leonardofibonacci/local/bin</string>
    <key>HOME</key>
    <string>/Users/leonardofibonacci</string>
    <key>CC_LOOP_TIMEOUT</key>
    <string>1800</string>
    <key>CC_LOOP_MAX_ERRORS</key>
    <string>5</string>
    <key>CC_LOOP_COOLDOWN</key>
    <string>300</string>
    <key>CC_LOOP_RATE_BACKOFF</key>
    <string>3600</string>
  </dict>
</dict>
</plist>
