# launchd/com.user.ecosystem-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.ecosystem-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/leonardofibonacci/code/scripts/ecosystem-health.sh &gt; /Users/leonardofibonacci/.claude/ecosystem-watchdog.last.json</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/leonardofibonacci</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>StandardOutPath</key>
  <string>/Users/leonardofibonacci/.claude/ecosystem-watchdog.stdout</string>

  <key>StandardErrorPath</key>
  <string>/Users/leonardofibonacci/.claude/ecosystem-watchdog.stderr</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/leonardofibonacci/local/bin:/Users/leonardofibonacci/local/opt/node/bin</string>
    <key>HOME</key>
    <string>/Users/leonardofibonacci</string>
  </dict>
</dict>
</plist>
